# chirpstack/chirpstack-gateway-mesh/files/uci-defaults/90_add_mesh_filters
uci get chirpstack-gateway-mesh.mesh_filters
if [ $? -eq 0 ]; then
  exit
fi

uci set chirpstack-gateway-mesh.mesh_filters='mesh_filters'
uci commit
